# Dockerfile/transcode/rootfs/usr/local/etc/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx

